# deploy/zabka-backend.service
# systemd unit for the zabka-dashboard backend
# Place at: /etc/systemd/system/zabka-backend.service
# After editing: sudo systemctl daemon-reload && sudo systemctl restart zabka-backend
#
# Changes vs the original single-worker unit:
#   - uvicorn called directly (not via python -m backend.main) so --workers works
#   - 1 worker on purpose (the app is not multi-worker-safe yet; see the note on
#     ExecStart below) - do NOT raise this without fixing per-process caches/locks
#   - access log off in production (reduces disk I/O under load)
#   - LimitNOFILE raised to match /etc/security/limits.conf

[Unit]
Description=Zabka Dashboard Backend
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=exec
User=zabka
Group=zabka
WorkingDirectory=/home/zabka/zabka-stats

# Environment - secrets live in /home/zabka/.env (chmod 600, not in repo)
EnvironmentFile=-/home/zabka/.env

# Core settings
Environment=ZABKA_DB=/home/zabka/zabka-stats/data/zabka.duckdb
Environment=REDIS_SOCKET=/run/redis/redis-server.sock

# Run uvicorn directly so --workers is honoured.
# ONE worker on purpose. The app is not multi-worker-safe yet: each worker keeps
# its own in-process caches and per-process stampede locks (so N workers = N cold
# caches and N-way cache stampedes), and the /api/snapshot -> ETL path opens a
# read-write DuckDB connection that would collide with other workers' read-only
# connections. Until those are made cross-process safe, stay at 1. A single
# worker is plenty for this traffic behind nginx's microcache + Redis.
ExecStart=/home/zabka/zabka-stats/venv/bin/uvicorn backend.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --no-access-log

Restart=on-failure
RestartSec=5s

# Match the system-wide ulimit set by harden_vps.sh
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
